# os/systemd-system/systemd-remount-fs.service.d/10-goblins-os-composefs.conf
[Unit]
# Fedora bootc/ostree composefs boots expose / as an overlay. systemd-remount-fs
# cannot reconfigure that mount; ostree-remount.service handles the real sysroot.
ConditionKernelCommandLine=!ostree
