# os/systemd/goblins-os-resident.service
[Unit]
Description=Goblins OS persistent Goblins AI runtime
After=goblins-os-core.service
Wants=goblins-os-core.service

[Service]
Type=simple
User=goblins-os
Group=goblins-os
EnvironmentFile=-/etc/goblins-os/environment
EnvironmentFile=-/etc/goblins-os/openai-secrets.env
Environment=GOBLINS_OS_RESIDENT_STATE=/var/lib/goblins-os/resident
Environment=GOBLINS_OS_RESIDENT_SOCKET=/run/goblins-os/resident.sock
ExecStart=/usr/libexec/goblins-os/goblins-os-resident
Restart=always
RestartSec=3s
StateDirectory=goblins-os/resident
RuntimeDirectory=goblins-os
RuntimeDirectoryMode=0750
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/goblins-os/resident
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
