# os/systemd/goblins-os-model-cache.service
[Unit]
Description=Goblins OS local model cache
Before=goblins-os-core.service

[Service]
Type=oneshot
User=goblins-os
Group=goblins-os
EnvironmentFile=-/etc/goblins-os/environment
StateDirectory=goblins-os/models goblins-os/installer goblins-os/session goblins-os/policy goblins-os/resident goblins-os/secrets
StateDirectoryMode=0750
ExecStart=/usr/bin/install -d -m 0700 /var/lib/goblins-os/secrets/openai
RemainAfterExit=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/goblins-os/installer /var/lib/goblins-os/session /var/lib/goblins-os/policy /var/lib/goblins-os/resident /var/lib/goblins-os/models /var/lib/goblins-os/secrets
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
